# deployment/systemd/ato-analyzer-worker.wsl-local.service
[Unit]
Description=ATO analyzer worker (WSL local — deterministic and model-assisted runs)
Documentation=file:///opt/ato-analyzer/docs/OPERATIONS_AND_RECOVERY.md
After=network-online.target postgresql.service ato-api.service
Wants=network-online.target
Requires=ato-api.service

[Service]
Type=simple
User=ato
Group=ato

WorkingDirectory=/opt/ato-analyzer

Environment=ATO_RUNTIME_CONFIG_PATH=/opt/ato-analyzer/runtime-config.json
Environment=PYTHONDONTWRITEBYTECODE=1

LoadCredential=database-dsn:/etc/ato-analyzer/credentials/database-dsn
LoadCredential=audit-hmac-key:/etc/ato-analyzer/credentials/audit-hmac-key
EnvironmentFile=-/etc/ato-analyzer/credentials/ato-local.env

ExecStart=/opt/ato-analyzer/venv/bin/ato-analyzer-worker

TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

Restart=on-failure
RestartSec=10

NoNewPrivileges=yes
CapabilityBoundingSet=
AmbientCapabilities=
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectKernelLogs=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
LockPersonality=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/ato-packages /opt/ato-analyzer/data/ato-storage

MemoryMax=2G
CPUQuota=200%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ato-analyzer-worker

[Install]
WantedBy=multi-user.target
